# nc_user_files_backup/rootfs/etc/s6-overlay/s6-rc.d/10-backup/run
#!/command/with-contenv bashio
# shellcheck shell=bash

RUNTIME_DIR="/run/nc_user_files_backup"
READY_FILE="${RUNTIME_DIR}/service.ready"

while [ ! -f "$READY_FILE" ]; do
    sleep 1
done

rm -f "${RUNTIME_DIR}/service.ready"
rm -f "${RUNTIME_DIR}/backup.lock"

# Start main application
exec /usr/bin/backup_scheduler
