# simple_dlna/rootfs/etc/s6-overlay/s6-rc.d/03-manager/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

RUNTIME_ENV="/run/simple_dlna/runtime.env"
WEB_APP="/usr/local/simple_dlna/web/filemanager.py"
READY_FILE="/run/simple_dlna/service.ready"

while [ ! -f "$READY_FILE" ]; do
    sleep 1
done

set -a
source "$RUNTIME_ENV"
set +a

exec python3 "$WEB_APP"
