# simple-dlna/rootfs/etc/s6-overlay/s6-rc.d/simple-dlna/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Backup Sync to USB addon starting"

# Start main application
exec /etc/simple-dlna/simple-dlna.sh
