# simple_dlna/rootfs/etc/s6-overlay/s6-rc.d/02-watcher/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

RUNTIME_ENV="/run/simple_dlna/runtime.env"
WATCHER_APP="/usr/local/simple_dlna/runtime/watcher.py"
READY_FILE="/run/simple_dlna/service.ready"

while [ ! -f "$READY_FILE" ]; do
    sleep 1
done

set -a
source "$RUNTIME_ENV"
set +a

exec python3 "$WATCHER_APP"
