# nc_user_files_backup/rootfs/etc/s6-overlay/s6-rc.d/15-cron/run
#!/command/with-contenv bashio
# shellcheck shell=bash

set -euo pipefail

bashio::log.green "Starting cron daemon"

exec crond -f -l 8 -L /var/log/cron.log -c /etc/crontabs