# backup_sync/rootfs/etc/s6-overlay/s6-rc.d/01-service/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Backup Sync to USB addon starting"

# Start main application
exec /usr/bin/backup_sync.sh
