# teletorrent/rootfs/etc/s6-overlay/s6-rc.d/teletorrent/finish
#!/usr/bin/with-contenv bash

echo "[teletorrent] stopped"
