# teletorrent/rootfs/etc/s6-overlay/s6-rc.d/teletorrent/run
#!/usr/bin/with-contenv bash

echo "[teletorrent] starting worker"

exec python3 /app/worker.py
