# teletorrent/rootfs/etc/s6-overlay/s6-rc.d/teletorrent/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.notice "========================================"
bashio::log.notice "=== STeleTorrent addon starting ==="
bashio::log.notice "Starting at: $(date '+%Y-%m-%d %H:%M:%S')"
bashio::log.notice "========================================"

bashio::log.info "Starting worker"

exec python3 /app/worker.py


