# teletorrent/rootfs/etc/s6-overlay/s6-rc.d/teletorrent/run
#!/command/with-contenv bashio
# shellcheck shell=bash

# Addon startup info
# Prints version, name and startup timestamp
ADDON_VERSION="$(bashio::addon.version)"
ADDON_NAME=$(bashio::addon.name)

bashio::log "============================================="
bashio::log.blue "=== ${ADDON_NAME} ==="
bashio::log.blue "=== Version:  ${ADDON_VERSION} ==="
bashio::log.green "Starting at: $(date '+%Y-%m-%d %H:%M:%S')"
bashio::log.info "Worker starting ..."

exec python3 /app/worker.py


