# teletorrent/rootfs/etc/s6-overlay/s6-rc.d/teletorrent/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log "========================================"
bashio::log.green "=== STeleTorrent addon starting ==="
bashio::log.green "Starting at: $(date '+%Y-%m-%d %H:%M:%S')"
bashio::log "========================================"

bashio::log.green "Starting worker"

exec python3 /app/worker.py


