# simple-dlna/rootfs/etc/s6-overlay/s6-rc.d/01-service/run
#!/command/with-contenv bashio
# shellcheck shell=bash

rm -f /run/simple-dlna/service.ready
exec /usr/bin/simple-dlna-service.sh