# backup_sync/rootfs/etc/services.d/backup_sync/finish
#!/command/with-contenv bashio
# shellcheck shell=bash

exit_code="$1"

# bashio::log.info "Exit code: $exit_code"

case "$exit_code" in

  256)
    bashio::log.info "Addon stopped"
    ;;

  1)
    bashio::log.fatal "Startup error — halting container"
    exec /run/s6/basedir/bin/halt
    ;;

  *)
    bashio::log.fatal "Addon crashed, halting container"
    exec /run/s6/basedir/bin/halt
    ;;

esac
