# caddy_proxy/rootfs/etc/s6-overlay/s6-rc.d/20-flask/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

bashio::log.green "Запуск Flask GUI"

cd /web || exit 1
exec gunicorn --bind 0.0.0.0:5000 app:app
