# caddy_proxy/rootfs/etc/s6-overlay/s6-rc.d/20-flask/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

cd /web || exit 1
exec python3 /web/app.py