# simple_dlna/rootfs/etc/s6-overlay/s6-rc.d/01-service/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

exit_code="$1"

BASE_DIR="/usr/local/simple_dlna"

emit() {
    python3 "${BASE_DIR}/ha/emit_cli.py" "$@" || true
}

if [ "$exit_code" -eq 0 ]; then
    emit service_state '{"reason":"stopping"}'
    bashio::log.green "Stopped at: $(date)"
    bashio::log.green "Service terminated normally. Good-bye!"
else
    bashio::log "========================================"
    bashio::log.red "Fail at: $(date)"
    emit service_state '{"reason":"crashed"}'
    bashio::log.red "Service crashed (exit code: ${exit_code}) — halting container"

    if [ -f /run/s6-linux-init-container-results/exitcode ]; then
        echo 0 > /run/s6-linux-init-container-results/exitcode
    fi

    exec /run/s6/basedir/bin/halt
fi

exit 0
