# simple_dlna/rootfs/etc/s6-overlay/s6-rc.d/01-service/run
#!/command/with-contenv bashio
# shellcheck shell=bash

rm -f /run/simple_dlna/service.ready
exec /usr/bin/simple_dlna_service.sh
