# torrent_bridge/rootfs/etc/s6-overlay/s6-rc.d/torrent_bridge/finish
#!/usr/bin/with-contenv bash

echo "[torrent_bridge] stopped"
