# torrent_bridge/rootfs/etc/s6-overlay/s6-rc.d/torrent_bridge/run
#!/usr/bin/with-contenv bash

echo "[torrent_bridge] starting worker"

exec python3 /app/worker.py
