# nc_user_files_backup/rootfs/etc/s6-overlay/s6-rc.d/20-web/run
#!/command/with-contenv bashio
# shellcheck shell=bash

set -euo pipefail

RUNTIME_ENV="/run/nc_user_files_backup/runtime.env"
WEB_APP="/usr/local/nc_user_files_backup/web/filemanager.py"
READY_FILE="/tmp/nc_user_files_backup/service.ready"

while [ ! -f "$READY_FILE" ]; do
    sleep 1
done

set -a
source "$RUNTIME_ENV"
set +a

exec python3 "$WEB_APP"
