# teletorrent/rootfs/etc/s6-overlay/s6-rc.d/00-init/up
#!/usr/bin/execlineb -P
with-contenv
bash /etc/s6-overlay/s6-rc.d/00-init/up.sh
