# teletorrent/rootfs/etc/s6-overlay/s6-rc.d/00-init/up
#!/usr/bin/with-contenv sh

exec /app/teletorrent/core/config.py
