# teletorrent/rootfs/etc/s6-overlay/s6-rc.d/00-init/up
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

cd /app

export ADDON_VERSION="$(bashio::addon.version)"
export ADDON_NAME="$(bashio::addon.name)"

python3 -m teletorrent.core.config
