# teletorrent/rootfs/etc/s6-overlay/s6-rc.d/10-teletorrent/run
#!/usr/bin/with-contenv sh

export PYTHONPATH=/app

exec python3 /app/worker.py


