# teletorrent/rootfs/etc/s6-overlay/s6-rc.d/10-teletorrent/run
#!/usr/bin/with-contenv sh

exec env PYTHONPATH=/app python3 /app/worker.py
