# teletorrent/rootfs/etc/s6-overlay/s6-rc.d/10-teletorrent/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Worker starting ..."

exec python3 /app/worker.py


