# nc_user_files_backup/rootfs/etc/s6-overlay/s6-rc.d/00-init/up
#!/command/with-contenv bashio
# shellcheck shell=bash

exec /usr/bin/backup_init