# tg2transmission/rootfs/etc/s6-overlay/s6-rc.d/tg2transmission/finish
#!/usr/bin/with-contenv bash

echo "[tg2transmission] stopped"
